# jbonjean/etc/service/socklog/run
#!/bin/sh
set -eu
. /usr/lib/runit-lib

export UID=$(id -u nobody)
exec socklog unix /dev/log
